# Weather Update/weather-update.timer
[Unit]
Description=Hourly Timer

[Timer]
OnCalendar=hourly
Persistent=true
AccuracySec=30s

[Install]
WantedBy=default.target
